# systemd/led-matrix-dashboard.service
[Unit]
Description=Raspberry Pi LED Matrix dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=human
Group=human
WorkingDirectory=/home/human/Raspberry-Pi-LED-Matrix
EnvironmentFile=-/home/human/Raspberry-Pi-LED-Matrix/.env
ExecStart=/usr/bin/python3 /home/human/Raspberry-Pi-LED-Matrix/dashboard_server.py --host 0.0.0.0
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
